# oma-deployment-package/configs/oma-autologin.service
[Unit]
Description=OMA Auto-login Setup Wizard
Documentation=https://github.com/DRDAVIDBANNER/X-Vire
After=multi-user.target network.target mariadb.service oma-api.service volume-daemon.service
Wants=network.target

[Service]
Type=idle
User=oma_admin
Group=oma_admin
TTY=/dev/tty1
ExecStart=/opt/ossea-migrate/oma-setup-wizard.sh
StandardInput=tty
StandardOutput=tty
StandardError=tty
Restart=no
RemainAfterExit=yes

# Environment variables
Environment=HOME=/home/oma_admin
Environment=USER=oma_admin
Environment=TERM=xterm-256color

# Security settings
NoNewPrivileges=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target






