# deployment/sha-appliance/configs/sendense-hub.service
[Unit]
Description=Sendense Hub Appliance (SHA) API Server
Documentation=https://docs.sendense.com/sha/
After=network.target mariadb.service volume-daemon.service
Wants=volume-daemon.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/sendense
ExecStart=/usr/local/bin/sendense-hub -port=8082 -db-host=localhost -db-port=3306 -db-name=migratekit_oma -db-user=oma_user -db-pass=oma_password
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sendense-hub

# Environment
Environment="PORT=8082"
Environment="LOG_LEVEL=info"
Environment="CONFIG_PATH=/etc/sendense/config.yaml"
Environment="DB_HOST=localhost"
Environment="DB_PORT=3306"
Environment="DB_NAME=migratekit_oma"
Environment="DB_USER=oma_user"
Environment="DB_PASS=oma_password"

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/sendense /var/lib/sendense /opt/sendense

[Install]
WantedBy=multi-user.target
